# etc/supervisor.conf
[program:bankomat]
command = /home/eli/bankomat/venv/bin/gunicorn api.app:app --bind unix:/tmp/bankomat.sock --log-file -
directory = /home/eli/bankomat
user = eli
