# etc/supervisor.conf
[program:bankomat]
command = /home/eli/bankomat/venv/bin/gunicorn webapp.webapp:app --log-file - --bind unix:/tmp/bankomat.sock
directory = /home/eli/bankomat
user = eli
